# backend/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/nginx/error.log
stdout_logfile=/var/log/nginx/access.log

[program:php-fpm]
command=php-fpm
autostart=true
autorestart=true
stderr_logfile=/var/log/php-fpm.log
stdout_logfile=/var/log/php-fpm.log

[program:laravel-setup]
command=/bin/bash -c "php artisan config:cache && php artisan route:cache && php artisan migrate --force"
autostart=true
autorestart=false
startsecs=0
stderr_logfile=/var/log/laravel-setup.log
stdout_logfile=/var/log/laravel-setup.log